# crontab
0 *12 * * * /home/pickle/Documents/my-vimrc/setup/backup.sh
